# aave-weth-isolation-monitor.service
[Unit]
Description=Aave Mantle WETH Isolated Debt Ceiling monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/aave-mantle-bot/aave_weth_isolation_monitor.py --interval 10 --chat-id 231918072 --min-free-usd 500 --cooldown 5
Restart=always
RestartSec=5
WorkingDirectory=/opt/aave-mantle-bot

[Install]
WantedBy=default.target
